# schedule.cron
0 * * * * /home/sonia/hackathons/log_rotator/rotator.sh
